# crontab
0 4 * * * (cd /home/mikew/projects/radio/cdbs/cdbs && make > crontab.log 2>&1)  

